# autostart/launch.sh.desktop
[Desktop Entry]
Exec=/home/kev/Downloads/launch.sh
Icon=dialog-scripts
Name=launch.sh
Path=
Type=Application
X-KDE-AutostartScript=true
